# src/samples/hello.ini
[program:echoserver]
command=/Library/Java/JavaVirtualMachines/jdk1.8.0_25.jdk/Contents/Home/bin/java -jar /data/giko/steve/sample/hello.jar
process_name=%(program_name)s
directory=/data/giko/steve/sample/
autostart=true
autorestart=true
stderr_logfile=/data/giko/steve/sample/logs/%(program_name)s.err.log
stdout_logfile=/data/giko/steve/sample/logs/%(program_name)s.out.log